# data/indicator-bluetooth.service.in
[Unit]
Description=Indicator Bluetooth Backend
PartOf=graphical-session.target
PartOf=indicators.target

[Service]
ExecStart=@pkglibexecdir@/indicator-bluetooth-service
Restart=on-failure

[Install]
WantedBy=indicators.target
